# supervisor/web.conf
[program:webserver]
directory=/root/hzj/BaseFastapiOrm/
command=/root/hzj/BaseFastapiOrm/env/bin/uvicorn main:app --host 0.0.0.0 --reload
autostart = true
